# templates/postgresql.service.j2
[Service]
Environment=PGDATA={{ postgres_server_path }}/
